# supervisord.conf
; Sample supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Note: shell expansion ("~" or "$HOME") is not supported.  Environment
; variables can be expanded using this syntax: "%(ENV_HOME)s".
[unix_http_server]
file=%(here)s/runtime/supervisor.sock

[supervisord]
logfile=%(here)s/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=%(here)s/runtime/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
directory=%(here)s

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://%(here)s/runtime/supervisor.sock

[program:flake]
command=python2.7 flake.py --logging=warning --worker_id=%(process_num)d --unix_socket=%(here)s/runtime/flake-%(process_num)02d.sock
process_name=%(program_name)s%(process_num)d
numprocs=4
numprocs_start=1
autostart=true
autorestart=true
stdout_logfile=%(here)s/log/flake.log
redirect_stderr=true
